# app/rootfs/etc/services.d/skygear-server/run
#!/usr/bin/with-contenv sh

/usr/local/bin/skygear-server /usr/src/app/config.ini
